# luci-app-pppoe-user/root/etc/init.d/pppoe-user
#!/bin/sh /etc/rc.common
# Copyright (C) 2018-2020 Lienol <lawlienol@gmail.com>

START=59

CONFIG="*"
CHAP_SECRETS=/var/etc/chap-secrets
TEMP=/tmp/$CONFIG.tmp

setup_users() {
	config_get enabled $1 enabled
	[ "$enabled" -eq 0 ] && return 0
	config_get username $1 username
	config_get password $1 password
	config_get ipaddress $1 ipaddress
	[ -n "$username" -a -n "$password" ] && echo "$username $CONFIG $password $ipaddress" >> $CHAP_SECRETS
}

del_user()
{
	cat $CHAP_SECRETS | grep -v $CONFIG > $TEMP
	cat $TEMP > $CHAP_SECRETS
	rm $TEMP
}

start() {
	config_load pppoe-user
	enabled=$(config_t_get service enabled)
	[ "$enabled" -eq 0 ]  && exit 0
	del_user
	config_foreach setup_users user
	ln -sfn $CHAP_SECRETS /etc/ppp/chap-secrets
}
