# luci-app-pppoe-user/root/etc/init.d/pppoe-user
#!/bin/sh /etc/rc.common

START=99

CHAP_SECRETS=/etc/ppp/chap-secrets

setup_user() {
	local cfg="$1"
	config_get enabled $1 enabled
	[ "$enabled" -eq 0 ] && return 0
	config_get username $1 username
	config_get servicename $1 servicename
	[ -n "${servicename}" ] || servicename="*"
	config_get password $1 password
	config_get ipaddress $1 ipaddress
	config_get macaddress $1 macaddress
	[ -n "$username" -a -n "$password" ] && echo "$username $servicename $password $ipaddress" >> $CHAP_SECRETS
	config_get qos $1 qos
	[ -n "${qos}" ] || qos="1"
	config_get download $1 download
	config_get upload $1 upload

	uci set nft-qos.$username=client >/dev/null 2>&1
	uci set nft-qos.$username.hostname=$username >/dev/null 2>&1
	uci set nft-qos.$username.macaddr=$macaddress >/dev/null 2>&1
	uci set nft-qos.$username.urate=$upload >/dev/null 2>&1
	uci set nft-qos.$username.drate=$download >/dev/null 2>&1
	uci commit nft-qos

        chmod 600 /etc/ppp/*-secrets
}

delete_user() {
	rm -rf $CHAP_SECRETS
}

start() {
	config_load pppoe-user
	delete_user
	config_foreach setup_user user
}

stop() {
	delete_user
	killall -q pppoe-user
}

