# luci-app-usermanager/root/etc/uci-defaults/80_usermanager
#!/bin/sh

touch /etc/config/usermanager
uci set usermanager.first=first
uci set usermanager.second=second
uci commit

return 0
