# luci-app-pppoeuser/root/etc/init.d/pppoeuser
#!/bin/sh /etc/rc.common

START=40
STOP=40

CURRENT_DATE=$(date +"%Y-%m-%d")
CHAP_SECRETS=/etc/ppp/chap-secrets

setup_user() {
    local cfg="$1"
    config_get enabled "$cfg" enabled
    config_get username "$cfg" username
    config_get servicename "$cfg" servicename
    [ -n "$servicename" ] || servicename="*"
    config_get password "$cfg" password
    config_get ipaddr "$cfg" ipaddr
    config_get macaddr "$cfg" macaddr
    config_get package "$cfg" package
    config_get qos "$cfg" qos
    config_get urate "$cfg" urate
    config_get drate "$cfg" drate
    config_get unit "$cfg" unit
    [ -n "$unit" ] || unit="kbytes"
    config_get connect "$cfg" connect
    config_get expires "$cfg" expires

    if [ "$expires" == "$CURRENT_DATE" ]; then
        uci set pppoeuser."$cfg".enabled=0
        uci commit pppoeuser
    fi

    # Do not append to the CHAP_SECRETS file if the user is not enabled
    if [ "$enabled" -eq 1 ] && [ "$expires" != "$CURRENT_DATE" ]; then
        echo "$username $servicename $password $ipaddr" >> "$CHAP_SECRETS"
    fi
}

daily_check() {
    (crontab -l ; echo "30 8 * * * /etc/init.d/pppoeuser start") | sort - | uniq - | crontab -
}

start() {
    cp -rf "$CHAP_SECRETS" /etc/ppp/user-chap-secrets
    > "$CHAP_SECRETS"
    daily_check
    config_load pppoeuser
    config_foreach setup_user user
    chmod 600 "$CHAP_SECRETS"
}

stop() {
    killall -q pppoeuser
}
